# systemd/user/mic-notify.service
[Unit]
Description=Microphone active notifications
After=graphical-session.target pipewire-pulse.service

[Service]
ExecStart=%h/.local/bin/mic-notify
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
