# systemd/user/rclone@.service
[Unit]
Description=Rclone mount for %I
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
Environment=RCLONE_CONFIG=%h/.config/rclone/rclone.conf

# Auto-create the specific mount folder
ExecStartPre=/usr/bin/mkdir -p %h/Cloud/%i

# The Mount Command
ExecStart=/usr/bin/rclone mount %i: %h/Cloud/%i \
    --vfs-cache-mode full \
    --vfs-cache-max-size 5G \
    --vfs-cache-max-age 168h \
    --vfs-cache-poll-interval 1m \
    --vfs-write-back 5s \
    --dir-cache-time 5m \
    --buffer-size 32M \
    --cache-dir=%h/.cache/rclone/%i \
    --timeout 1h \
    --contimeout 60s \
    --retries 3 \
    --low-level-retries 10 \
    --no-modtime \
    --no-checksum \
    --log-level INFO \
    --log-file %h/.cache/rclone/logs/rclone-%i.log

# Lazy Unmount (Detaches immediately, prevents shutdown hangs)
ExecStop=/bin/fusermount3 -uz %h/Cloud/%i

Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical-session.target
