# systemd/user/study-library-sync.service
[Unit]
Description=Sync study library from Google Drive
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone sync gdrive:FAU/Library %h/.local/share/study-library --progress

[Install]
WantedBy=graphical-session.target
